# root/etc/udev/rules.d/99-displays.rules
SUBSYSTEM=="drm", ACTION=="change", DEVTYPE= RUN+="/usr/bin/displays-rearrange"
